# uwsgi.conf
start on runlevel [2345]
stop on runlevel [!2345]

chdir /var/web
exec /var/web/env/bin/uwsgi uwsgi.ini
